# docker/crontab
0 13,21 * * * /usr/local/bin/python -m hh_applicant_tool -vv -c /app/config.json apply-similar
0 19 * * * /usr/local/bin/python -m hh_applicant_tool -vv -c /app/config.json refresh-token
